# rootfs/etc/systemd/system/apps-bundle.target
[Unit]
Description=Apps Bundle
Requires=basic.target
After=basic.target
