# rootfs/etc/systemd/system/apps-bundle.target
[Unit]
Description=Apps Bundle
PartOf=apps.service
